# apps/pos-server/pos-server.service
[Unit]
Description=denkfink POS Thermal Printer Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=skander
WorkingDirectory=/home/skander/meinungeheuer/apps/pos-server
ExecStart=/usr/bin/python3 print_server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
